# files/containers/monitoring/otel-collector/otel-collector.container
[Unit]
Description=OpenTelemetry Collector (quadlet)

[Container]
Image=docker.io/otel/opentelemetry-collector-contrib:latest
ContainerName=otel-collector

#flag for autoupdates
Label=io.containers.autoupdate=registry

# Run as local user to allow access to data directories
User=$(id -u):$(id -g)

Network=monitoring.network

PublishPort=4317:4317 # OTLP gRPC
PublishPort=4318:4318 # OTLP HTTP
PublishPort=8888:8888 # Metrics endpoint

Volume=/etc/localtime:/etc/localtime:ro
Volume=/mnt/data/otel-collector/config:/etc/otelcol-contrib
Volume=/mnt/data/otel-collector/data:/data

# Configure environment variables to direct telemetry to proper backends
Environment=OTEL_EXPORTER_OTLP_METRICS_ENDPOINT=http://mimir:9009/otlp
Environment=OTEL_EXPORTER_OTLP_LOGS_ENDPOINT=http://loki:3100/loki/api/v1/push
Environment=OTEL_EXPORTER_OTLP_TRACES_ENDPOINT=http://tempo:4317

# Health check
HealthCmd=wget -q --spider http://localhost:13133 || exit 1
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=60s
HealthTimeout=10s

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target