# config/systemd/system/mnt-media.mount
[Unit]
Description=Media mount
After=network.target

[Mount]
What=10.1.0.24:/Media
Where=/var/mnt/media
Type=nfs
Options=_netdev,auto

[Install]
WantedBy=multi-user.target
