# files/containers/monitoring/tempo/tempo.container
[Unit]
Description=Tempo Distributed Tracing (quadlet)

[Container]
Image=docker.io/grafana/tempo:latest
ContainerName=tempo

#flag for autoupdates
Label=io.containers.autoupdate=registry

# Run as local user to allow access to data directories
User=0

Network=monitoring.network

Volume=/etc/localtime:/etc/localtime:ro
Volume=/mnt/data/tempo/config:/etc/tempo
Volume=/mnt/data/tempo/data:/data

Environment=TEMPO_STORAGE_TRACE_RETENTION=2160h
Environment=TEMPO_SERVER_HTTP_LISTEN_PORT=3200

# Health check
HealthCmd=wget -q --spider http://localhost:3200/ready || exit 1
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=60s
HealthTimeout=10s

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target